# etc/slurmbrokerd.service.in
[Unit]
Description=Slurm cross-domain broker daemon
After=network-online.target remote-fs.target munge.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/slurmbrokerd
EnvironmentFile=-/etc/default/slurmbrokerd
User=slurm
Group=slurm
ExecStart=@sbindir@/slurmbrokerd $SLURMBROKERD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
LimitNOFILE=65536

# Uncomment the following lines to disable logging through journald.
# NOTE: It may be preferable to set these through an override file instead.
#StandardOutput=null
#StandardError=null

[Install]
WantedBy=multi-user.target
